# supervisord.conf
[supervisord]
nodaemon=true

[program:glusterd]
command=service glusterd start
autorestart=true

[program:rpcbind]
command=/usr/sbin/rpcbind -f
