# gluster-object/CentOS/docker-gluster-swift/memcached.service
[Unit]
Description=Memcached Service

[Service]
ExecStart=/usr/bin/memcached -u root
Restart=on-abort

[Install]
WantedBy=multi-user.target
